# deploy/docker-web/supervisor.conf
[supervisord]
nodaemon=true

[program:django]
command=/bin/bash -c "cd /home/user/code/backend/sito && python manage.py runserver 0.0.0.0:8000"
